# src/etc/rc.local
mkdir -p /var/log/nginx
chmod a+rw /var/log/nginx
exit 0

